# templates/default/etc/systemd/system/elrond-node@.service.erb
[Install]
WantedBy=multi-user.target

[Unit]
Description=Elrond Node %i
After=network-online.target

[Service]
<% if @arwen -%>
; node is looking for this in CWD and it ignores PATH
Environment="ARWEN_PATH=/opt/elrond/bin/arwen"
<% end -%>
User=elrond-node-%i
; allow parameters for template unit
EnvironmentFile=<%= @var_dir %>/node-%i/config/service.env
; technically, this also holds config, but the node binary makes it difficult
; to separate config in etc and variable data in var
WorkingDirectory=/opt/var/elrond/node-%i
ExecStart=/opt/elrond/bin/node -use-log-view -log-level ${LOG_LEVEL} -rest-api-interface localhost:${REST_API_PORT}
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=3
LimitNOFILE=4096
